# controller/src/main/resources/installers/com.edgeguardian.agent.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.edgeguardian.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/edgeguardian-agent</string>
        <string>-config</string>
        <string>/etc/edgeguardian/agent.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/var/lib/edgeguardian</string>
    <key>StandardOutPath</key>
    <string>/var/log/edgeguardian-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/edgeguardian-agent.log</string>
</dict>
</plist>
